# deploy/linux/leba-acme-renew.service
[Unit]
Description=Leba ACME certificate renew (lego via admin API)
After=network-online.target leba.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/leba/leba.env
# Expect LEBA_ADMIN_AUTH=user:pass and optional LEBA_ADMIN_URL
ExecStart=/usr/bin/curl -fsS -u "${LEBA_ADMIN_AUTH}" -X POST \
  "${LEBA_ADMIN_URL:-http://127.0.0.1:8404}/admin/certificates/renew"
# tls-reload is triggered by the renew API; optional explicit reload:
# ExecStartPost=/usr/bin/curl -fsS -u "${LEBA_ADMIN_AUTH}" -X POST \
#   "${LEBA_ADMIN_URL:-http://127.0.0.1:8404}/admin/tls-reload"

[Install]
WantedBy=multi-user.target
